# com.elitephysio.marketing.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the marketing / NPS poller. Runs marketing_poll.sh every
  10 minutes (StartInterval 600 seconds).

  DO NOT load this until the marketing system is ready for SAFE_MODE testing
  (Resend / Tally / sheet configured). Until then the poller is harmless, but
  there is no reason to run it.

  Install:
    cp com.elitephysio.marketing.poll.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.elitephysio.marketing.poll.plist
  Pause:
    launchctl unload ~/Library/LaunchAgents/com.elitephysio.marketing.poll.plist
  Status:
    launchctl list | grep marketing
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elitephysio.marketing.poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/martinloughran/cliniko-dropoffs/marketing_poll.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/marketing-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/marketing-launchd.log</string>
</dict>
</plist>
