# com.elitephysio.progress.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the weekly off-track / progress review.
  Runs progress_poll.sh every Monday at 07:30 — it scans the previous
  completed Mon-Sun week and refreshes the "Off-Track Review" tab.

  Install:
    cp com.elitephysio.progress.weekly.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.elitephysio.progress.weekly.plist
  Pause:
    launchctl unload ~/Library/LaunchAgents/com.elitephysio.progress.weekly.plist
  Status:
    launchctl list | grep progress
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elitephysio.progress.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/martinloughran/cliniko-dropoffs/progress_poll.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/progress-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/progress-launchd.log</string>
</dict>
</plist>
