# com.elitephysio.eod.oneshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ONE-SHOT EOD handover post — fires at 20:15 tonight (2026-05-19) on request.
  eod_oneshot.sh deletes this job after it runs, so it never repeats.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elitephysio.eod.oneshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/martinloughran/cliniko-dropoffs/eod_oneshot.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/eod-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/eod-launchd.log</string>
</dict>
</plist>
