# com.elitephysio.dropoff.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the daily patient drop-off automation.
  Runs run_daily.sh every morning at 07:00 — pulls yesterday's
  cancellations / DNAs / no-rebooks from Cliniko (rolling 10-day rescan),
  AI-categorises body area, writes the Google Sheet, refreshes all dashboard
  tabs, and sends the Slack DMs (physios + reception + Ops Manager summary).

  run_daily.sh retries once after 120s on failure (covers a Wi-Fi/DNS blip
  just after the Mac wakes). If the Mac is asleep at 07:00, launchd fires the
  job as soon as it next wakes.

  Install:
    cp com.elitephysio.dropoff.daily.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.elitephysio.dropoff.daily.plist
  Pause:
    launchctl unload ~/Library/LaunchAgents/com.elitephysio.dropoff.daily.plist
  Status:
    launchctl list | grep elitephysio
  Trigger manually:
    launchctl start com.elitephysio.dropoff.daily
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elitephysio.dropoff.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/martinloughran/cliniko-dropoffs/run_daily.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/launchd-stderr.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/launchd-stderr.log</string>
</dict>
</plist>
