# com.elitephysio.eod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the end-of-day stats report.
  Runs eod_poll.sh 15 minutes before each admin shift ends — it posts the
  stats table to the #eod-claude Slack channel for the handover email.

    Shifts:  Mon-Thu 13:00 / 16:30 / 21:00   ->  report 12:45 / 16:15 / 20:45
             Fri     16:00                   ->  report 15:45

  launchd weekday codes: 1=Mon 2=Tue 3=Wed 4=Thu 5=Fri.

  Install:
    cp com.elitephysio.eod.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.elitephysio.eod.plist
  Pause:
    launchctl unload ~/Library/LaunchAgents/com.elitephysio.eod.plist
  Status:
    launchctl list | grep eod
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elitephysio.eod</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/martinloughran/cliniko-dropoffs/eod_poll.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- Mon-Thu 12:45 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>45</integer></dict>
        <!-- Mon-Thu 16:15 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>15</integer></dict>
        <!-- Mon-Thu 20:45 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>45</integer></dict>
        <!-- Fri 15:45 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/eod-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinloughran/cliniko-dropoffs/logs/eod-launchd.log</string>
</dict>
</plist>
